# debian/postrm
#!/bin/sh
#

if [ "$1" = purge ]; then
    update-rc.d arpsponge remove
fi
